# hypr/clipboard-notify.service
[Unit]
Description=Clipboard notification service
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/rolle/.local/bin/clipboard-notify.sh
Restart=always
RestartSec=5
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-1
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus

[Install]
WantedBy=default.target
